# files/logstash-forwarder.ini
[program:logstash-forwarder]
command = /opt/logstash-forwarder/bin/logstash-forwarder -config /etc/logstash-forwarder.json
